# vps1/services/rootless/blocky/tailscale-hosts/blocky-tailscale.container.jinja
[Unit]
Description=Blocky Headscale hosts-file adapter

[Container]
ContainerName=blocky-tailscale
Image=docker.io/oven/bun:alpine
Pull=newer
Environment=BASE_DOMAIN={{ base_domain }}
PublishPort=127.0.0.1:16932:16932
Volume=.:/src:ro,Z
Volume=/home/core/containers-git/vps1/volumes/blocky-tailscale-node:/app/node_modules:Z
Volume=/home/core/containers-git/vps1/volumes/headscale-data:/headscale:ro,z
WorkingDir=/app
Exec=sh -lc "cp -a /src/. /app && bun install && bun run start"

[Service]
TimeoutStartSec=15min
ExecStartPre=/usr/bin/mkdir -p /home/core/containers-git/vps1/volumes/blocky-tailscale-node

[Install]
WantedBy=default.target
