# vps1/services/rootless/rybbit/rybbit-client.container
[Unit]
Description=Rybbit web client
Requires=rybbit-backend.container
After=rybbit-backend.container

[Container]
ContainerName=rybbit-client
Image=ghcr.io/rybbit-io/rybbit-client:v2.7.0
Pull=newer
PublishPort=127.0.0.1:52363:3002
EnvironmentFile=./env/client.env
Environment=NODE_ENV=production

[Service]
TimeoutStartSec=15min
# Next.js exits with status 143 when stopped by podman/systemd
SuccessExitStatus=143

[Install]
WantedBy=default.target
