# vps1/services/matrix-sender/matrix-sender.container.jinja
[Unit]
Description=Matrix sender API
Wants=tuwunel.service caddy.service
After=tuwunel.service caddy.service

[Container]
ContainerName=matrix-sender
Image=docker.io/oven/bun:debian
Pull=newer
AddHost=chat.{{ base_domain }}:host-gateway
PublishPort=127.0.0.1:52498:8080
Environment="MATRIX_HOMESERVER=https://chat.{{ base_domain }}"
Environment="MATRIX_USERNAME={{ matrix_username }}"
Environment="MATRIX_PASSWORD={{ matrix_password }}"
Environment="MATRIX_RECOVERY_KEY={{ matrix_recovery_key }}"
Volume=.:/src:ro,Z
Volume=/home/core/containers-git/vps1/volumes/matrix-sender-node:/app/node_modules:Z
Volume=/home/core/containers-git/vps1/volumes/matrix-sender-data:/app/data:Z
WorkingDir=/app
Exec=sh -lc "cp -a /src/. /app && bun install && bun run start"

[Service]
TimeoutStartSec=15min
ExecStartPre=/usr/bin/mkdir -p /home/core/containers-git/vps1/volumes/matrix-sender-node
ExecStartPre=/usr/bin/mkdir -p /home/core/containers-git/vps1/volumes/matrix-sender-data

[Install]
WantedBy=default.target
